# root/etc/services.d/rsyncd/run
#!/usr/bin/with-contenv bash

echo "Starting Rsync Daemon."
exec /usr/bin/rsync --no-detach --daemon --config /etc/rsyncd.conf
